# deploy/pa-fetcher@.service
# Systemd template service for Personal Assistant Fetcher Workers
#
# This is a template unit that allows running multiple worker instances.
# The %i placeholder is replaced with the instance number.
#
# Installation:
# 1. Copy this file to /etc/systemd/system/pa-fetcher@.service
# 2. Run: sudo systemctl daemon-reload
# 3. Start workers: sudo systemctl start pa-fetcher@{1..3}
# 4. Enable on boot: sudo systemctl enable pa-fetcher@{1..3}
#
# Or use the manager script:
#   pa-worker-manager start   # Starts N workers based on fetch_worker_count setting
#   pa-worker-manager stop    # Stops all workers
#   pa-worker-manager status  # Shows status of all workers
#
# Commands for individual workers:
# - sudo systemctl status pa-fetcher@1     # Check worker 1 status
# - sudo systemctl restart pa-fetcher@2    # Restart worker 2
# - sudo journalctl -u pa-fetcher@1 -f     # Follow worker 1 logs
# - sudo journalctl -u 'pa-fetcher@*' -f   # Follow all worker logs

[Unit]
Description=Personal Assistant Fetcher Worker %i
Documentation=https://github.com/tolstuun/personal_assistant
After=network.target postgresql.service
# Allow multiple instances to start in parallel
PartOf=pa-fetcher.target

[Service]
Type=simple
User=assistant
Group=assistant

# Working directory (adjust path as needed)
WorkingDirectory=/opt/personal-assistant

# Python virtual environment
Environment="PATH=/opt/personal-assistant/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Application environment
Environment="PYTHONUNBUFFERED=1"
# Instance identifier available as WORKER_INSTANCE
Environment="WORKER_INSTANCE=%i"

# The worker command
ExecStart=/opt/personal-assistant/venv/bin/python -m src.workers.security_digest_worker

# Restart policy
Restart=always
RestartSec=10

# Process limits
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pa-fetcher-%i

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/
ReadWritePaths=/opt/personal-assistant/logs

[Install]
WantedBy=multi-user.target pa-fetcher.target
